# njr.conf
[program:njr]
command = nodejs /home/rtorrent/nodejs-rtorrent/index.js
directory = /home/rtorrent/nodejs-rtorrent
user = rtorrent
autostart = true
autorestart = true
startsecs=10